# deploy/eryndor.service
[Unit]
Description=Eryndor MMO Game Server
After=network.target
Documentation=https://github.com/jbuehler23/eryndor-mmo

[Service]
Type=simple
User=eryndor
Group=eryndor
WorkingDirectory=/opt/eryndor
ExecStart=/opt/eryndor/server
Restart=always
RestartSec=10

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/eryndor

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G

# Environment variables (loaded from /etc/eryndor/.env)
EnvironmentFile=/etc/eryndor/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eryndor

[Install]
WantedBy=multi-user.target
